# Linux CRON Jobs/crontab
0 3 * * * /usr/bin/php /var/www/delta/update_data.php
